# packaging/marathon.systemd.postinst
#!/bin/sh
set -e
systemctl enable marathon
